# config/launchd/com.tmux.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.tmux.server</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/tmux</string>
		<string>start-server</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>

	<key>ProcessType</key>
	<string>Background</string>

	<key>StandardOutPath</key>
	<string>/tmp/tmux-launchd-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/tmux-launchd-stderr.log</string>
</dict>
</plist>
